# configs/rtmpstats.cron
0,15,30,45 * * * * /opt/rtmpdash/stats.sh
